# ansible/roles/rrdstats/templates/generate-rrdtool-stats.timer.j2
# {{ ansible_managed }}

[Unit]
Description = Generate rrdtool stats minutely

[Timer]
OnCalendar = minutely

[Install]
WantedBy = timers.target

# vim: ft=dosini.j2:
