# ansible/roles/lke-nftables-update/templates/lke-nftables-update.timer.j2
[Unit]
Description=Hourly update of LKE IP whitelist in nftables

[Timer]
OnCalendar=hourly
RandomizedDelaySec=600
Persistent=true

[Install]
WantedBy=timers.target
