# ansible/roles/lke-nftables-update/templates/lke-nftables-update.service.j2
[Unit]
Description=Update LKE IP whitelist in nftables
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-lke-ips.sh
